# templates/geminabox.conf.erb
description "geminabox"
version "1.0"
author "David Goodlad"
 
env LANG=en_US.UTF-8
env APP_HOME=<%= root %>

respawn
start on runlevel [23]
 
exec start-stop-daemon --start --make-pidfile --pidfile /var/run/geminabox.pid --chuid <%= user %> --chdir <%= root %> --exec /usr/local/bin/thin -- start -e production -p <%= port %> >> <%= root %>/geminabox.log 2>&1
